# strange-music.service
[Unit]
Description=Play strange music for 5 minutes

[Service]
ExecStart=/usr/bin/strange-music 300
