# strange-music.timer
[Unit]
Description=Play strange music everyday

[Timer]
OnCalendar=daily
AccuracySec=24h
Persistent=true

[Install]
WantedBy=timers.target
